# services/getty@tty2.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/watch -n 5.0 -t /root/watch.sh
#StandardInput=tty
StandardOutput=tty
